# dist/rootfs/usr/lib/systemd/system/openlan-ceci@.service
[Unit]
Description=OpenLAN Ceci daemon for %I
After=network.target

[Service]
Type=notify
LimitNOFILE=102400
ExecStart=/usr/bin/openlan-ceci $OPTIONS -conf /etc/openlan/switch/ceci/%i.yaml
Restart=always

[Install]
WantedBy=multi-user.target
