# dist/rootfs/usr/lib/systemd/system/openceci@.service
[Unit]
Description=OpenCeci daemon for %I
After=network.target

[Service]
Type=notify
LimitNOFILE=102400
ExecStart=/usr/bin/openceci $OPTIONS -conf /etc/openlan/switch/ceci/%i.yaml
Restart=always

[Install]
WantedBy=multi-user.target
