# dist/rootfs/usr/lib/systemd/system/openceci@.service
[Unit]
Description=Open Cecilia daemon for %I
After=network.target

[Service]
Type=notify
LimitNOFILE=102400
EnvironmentFile=/etc/sysconfig/openlan/ceci.cfg
ExecStart=/usr/bin/openceci $OPTIONS -conf /etc/openlan/switch/ceci/%i.yaml
Restart=always

[Install]
WantedBy=multi-user.target
