# dist/rootfs/usr/lib/systemd/system/openlan-access@.service
[Unit]
Description=OpenLAN Access daemon for %I
After=network.target

[Service]
Type=notify
EnvironmentFile=/etc/sysconfig/openlan/access.cfg
ExecStart=/usr/bin/openlan-access $OPTIONS -conf /etc/openlan/%i.json
Restart=always

[Install]
WantedBy=multi-user.target
