# dist/rootfs/usr/lib/systemd/system/openlan-switch.service
[Unit]
Description=OpenLAN Central Switch daemon
After=network.target

[Service]
Type=once
EnvironmentFile=/etc/sysconfig/openlan/switch.cfg
ExecStart=/var/openlan/script/switch.sh start $OPTIONS
LimitNOFILE=102400

[Install]
WantedBy=multi-user.target
