# dist/rootfs/usr/lib/systemd/system/openlan-proxy.service
[Unit]
Description=OpenLAN Proxy daemon
After=network.target

[Service]
Type=notify
EnvironmentFile=/etc/sysconfig/openlan/proxy.cfg
ExecStart=/usr/bin/openlan-proxy $OPTIONS
LimitNOFILE=102400
Restart=always

[Install]
WantedBy=multi-user.target
