# dist/rootfs/usr/lib/systemd/system/openlan-switch.service
[Unit]
Description=OpenLAN Central Switch daemon
After=network.target
Requires=openlan-confd.service

[Service]
Type=notify
EnvironmentFile=/etc/sysconfig/openlan/switch.cfg
ExecStart=/var/openlan/script/switch.sh start $OPTIONS
LimitNOFILE=102400
Restart=always

[Install]
WantedBy=multi-user.target
